# systemd/clamav-update.timer
[Unit]
Description=Update clamav virus definitions weekly

[Timer]
OnCalendar=Mon, 01:00
Persistent=true

[Install]
WantedBy=timers.target
