# systemd/clamav-update.service
[Unit]
Description=ClamAV virus definition update

[Service]
Type=oneshot
ExecStart=/usr/bin/freshclam
