# systemd/reflector.timer
[Unit]
Description=Run reflector weekly

[Timer]
OnCalendar=Tue, 01:00
Persistent=true

[Install]
WantedBy=timers.target
