# systemd/restic-prune.timer
[Unit]
Description=Prune data from the restic repository every Friday

[Timer]
OnCalendar=Friday
Persistent=true

[Install]
WantedBy=timers.target
